# kde/.config/autostart/kde.config.desktop
[Desktop Entry]
Exec=/home/mid-kid/.config/kde.config
Icon=application-x-shellscript
Name=kde.config
Type=Application
X-KDE-AutostartScript=true
